# step-04/root-files/opt/flownative/supervisor/etc/conf.d/php-fpm.conf
[program:php-fpm]
process_name=%(program_name)s
command="%(ENV_PHP_BASE_PATH)s/sbin/php-fpm"
autostart=true
autorestart=true
numprocs=1
